# shipper/rpm/usr/lib/systemd/system/shipper.service
[Service]
ExecStart=/bin/sh -c '${JAVA} ${JAVA_ARGS} -cp "lib:plugins/*" -jar shipper-fat.jar'
User=shipper
Type=simple
WorkingDirectory=/usr/share/shipper
Environment=JAVA=/usr/lib/jvm/adoptopenjdk-15-hotspot-jre/bin/java
EnvironmentFile=-/etc/sysconfig/shipper

[Install]
WantedBy=multi-user.target


